# alpine-pkg/dockerd/etc/init.d/dockerd
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/bin/dockerd
command_args_foreground="--containerd /run/containerd/containerd.sock --default-runtime='oci' --add-runtime 'oci=/usr/bin/runc'"
pidfile="/run/$RC_SVCNAME/$RC_SVCNAME.pid"
name=$RC_SVCNAME
description="Dockerd"

start_pre()
{
    # Ensure Service DIR exits
    mkdir -p /run/$RC_SVCNAME/
    export DOCKER_RAMDISK=true
}
