# pkg/init/etc/inittab
# /etc/inittab

::sysinit:/etc/init.d/rcS
#::once:/etc/init.d/containerd
#::once:/etc/init.d/containers

::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

# Stuff to do for the 3-finger salute
::ctrlaltdel:/sbin/reboot

# Stuff to do before rebooting
::shutdown:/usr/sbin/killall5 -15
::shutdown:/bin/sleep 5
::shutdown:/usr/sbin/killall5 -9
::shutdown:/bin/echo "Unmounting filesystems"
::shutdown:/bin/umount -a -r

# Put a getty on the serial port
ttyS0::respawn:/sbin/agetty -a root -L ttyS0 115200 vt100

