# debian-pkg/openrc/etc/init.d/containerd
#!/sbin/openrc-run

supervisor=supervise-daemon
command=/usr/bin/containerd
command_background="yes"
pidfile="/run/$RC_SVCNAME/$RC_SVCNAME.pid"
name=$RC_SVCNAME
description="Containerd"

start_pre()
{
    # Ensure Service DIR exits
    mkdir -p /run/$RC_SVCNAME/
}
